# build/package/quad-ops.service
[Unit]
Description=Quadlet Operations Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/quad-ops
Restart=always
RestartSec=10
WorkingDirectory=/etc/quad-ops
Environment=CONFIG_PATH=/etc/quad-ops/config.yaml

[Install]
WantedBy=multi-user.target
