# build/package/quad-ops.service
[Unit]
Description=Quad-Ops Container Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/quad-ops sync
ExecStart=/usr/local/bin/quad-ops up
