# build/package/quad-ops.service
[Unit]
Description=Quad-Ops Container Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/quad-ops --daemon --interval 300
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
