# build/package/quad-ops-user.service
[Unit]
Description=Quad-Ops Container Manager (User)
After=network.target

[Service]
Type=simple
Environment=QUADLET_DIR=%h/.config/containers/systemd
Environment=MANIFESTS_DIR=%h/.config/quad-ops/manifests
ExecStart=/usr/local/bin/quad-ops
Restart=on-failure

[Install]
WantedBy=default.target
