# build/package/quad-ops-user.service
[Unit]
Description=Quad-Ops Container Manager (User)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/quad-ops sync --daemon --user
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
