# build/package/quad-ops.service
[Unit]
Description=Quad-Ops Container Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/quad-ops sync --daemon
ExecStartPost=-/usr/local/bin/quad-ops up
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30
WatchdogSec=60

[Install]
WantedBy=multi-user.target default.target
